# sipp01-scenario-Adv.timer
[Unit]
Description=Run sipp01-scenario-Adv.service every 35 minutes

[Timer]
OnCalendar=*:0/15
